# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=/usr/bin/mongod -D

[program:npm]
command=npm install --production
command=npm start
autostart=true
autorestart=true

